# printer-init/S81mjpg_streamer
#!/bin/sh /etc/rc.common
# Custom MJPEG streamer for local camera access
# Replaces Creality's crash-prone uvc_stream with a stable mjpg-streamer build
# Runs slightly after the camera device check (S80camera) at boot

START=81
USE_PROCD=1
STOP=01

STREAMER_DIR=/mnt/UDISK
STREAMER_BIN=$STREAMER_DIR/mjpg_streamer

streamer_status=off

start_service() {
        while :
        do
                if [ -e /dev/video0 ]; then
                        if [ "${streamer_status}" == off ]; then
                                time=$(date "+%Y-%m-%d %H:%M:%S")
                                echo -e "[ ${time} ]: starting custom mjpg_streamer" > /dev/ttyS0
                                cd $STREAMER_DIR && ./mjpg_streamer -i "input_uvc.so -d /dev/video0 -r 640x480 -f 25" -o "output_http.so -p 8081 -w ./www" &
                                streamer_status=on
                        fi
                else
                        if [ "${streamer_status}" == on ]; then
                                PIDS=`ps | grep mjpg_streamer | grep -v grep | awk '{print $1}'`
                                if [ -n "${PIDS}" ]; then
                                        kill -9 ${PIDS}
                                fi
                                streamer_status=off
                        fi
                fi
                sleep 0.5
        done
}
